# clamav/clamav-weekly-scan.timer
[Unit]
Description=ClamAV weekly full-system scan timer

[Timer]
OnCalendar=Sat *-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
